# base/.config/systemd/user/export-ssh-keys.service
[Unit]
Description=Update SSH keys

[Service]
Type=oneshot
ExecStart=%h/bin/export-ssh-keys

[Install]
WantedBy=default.target
