# base/.config/systemd/user/export-ssh-keys.timer
[Unit]
Description=Run export-ssh-keys.service every 30 minutes

[Timer]
OnBootSec=30min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
