# base/.config/systemd/user/pika-ssh-sync.service
[Unit]
Description=Sync SSH keys

[Service]
Type=oneshot
ExecStart=%h/.local/opt/pika-ssh-sync/pika-ssh-sync

[Install]
WantedBy=default.target
